# setup/launchd/com.fcu.parser.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fcu.parser</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>git -C ~/fcu pull --quiet &amp;&amp; cd ~/fcu/bid-scanner &amp;&amp; python parser.py --auto</string>
    </array>

    <!-- Daily at 6:30am and 1:00pm (catches portals that upload docs late morning) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/fcu/bid-scanner</string>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/fcu/bid-scanner/logs/parser.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/fcu/bid-scanner/logs/parser.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
